# systemd/samplebox.service
[Unit]
Description=samplebox

[Service]
ExecStart=/home/pi/samplebox/box.rb

[Install]
WantedBy=multi-user.target